# crontab
0 0 * * * usr/local/bin/python /app/src/script.py
