# crontab
0 0 * * * /usr/bin/python /app/src/script.py
